# templates/default/systemd/mysqld.service.erb
[Service]
User=<%= @config.run_user %>
Group=<%= @config.run_group %>
