# templates/default/systemd/mysqld.service.erb
[Unit]
Description=mysql_service[mysql-<%= @config.instance %>]
After=syslog.target
After=network.target

[Service]
Type=simple
User=<%= @config.run_user %>
Group=<%= @config.run_group %>
PermissionsStartOnly=true
ExecStartPre=<%= @mysql_systemd_start_pre %>

ExecStart=<%= @mysqld_bin %> --defaults-file=<%= @etc_dir %>/my.cnf --basedir=<%= @base_dir %>
ExecStartPost=<%= @mysql_systemd %>
TimeoutSec=300
Restart=on-failure
RuntimeDirectory=mysqld
RuntimeDirectoryMode=755

[Install]
WantedBy=multi-user.target
